# frozbot.service
[Unit]
Description=FrozBot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/frozbot
EnvironmentFile=-/path/to/your/frozbot/.env
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/path/to/your/frozbot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/path/to/your/frozbot/.venv/bin/python /path/to/your/frozbot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
