# frozbot.service
[Unit]
Description=FrozBot Discord Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/frozbot
Environment=PATH=/path/to/your/frozbot/.venv/bin
ExecStart=/path/to/your/frozbot/.venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Log rotation settings
LogRotateMaxBytes=1M
LogRotateMaxFiles=1
# Limit journal entries to prevent excessive logging
LogRateLimitIntervalSec=1
LogRateLimitBurst=1000

[Install]
WantedBy=multi-user.target
